# deploy/systemd/streamlit-app.service
[Unit]
Description=AI Call Intelligence Streamlit Application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/AI-Call-Intelligence
Environment="PATH=/home/ubuntu/AI-Call-Intelligence/venv/bin"
ExecStart=/home/ubuntu/AI-Call-Intelligence/venv/bin/streamlit run app.py --server.port 8501 --server.address 0.0.0.0 --server.headless true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
